# launchd/local.theme-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.theme-watcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>__HOME__/scripts/apply_theme.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>1</integer>

  <key>WatchPaths</key>
  <array>
    <string>__HOME__/Library/Preferences/.GlobalPreferences.plist</string>
  </array>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/theme-watcher.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/theme-watcher.log</string>
</dict>
</plist>


